# deploy/systemd/ros2-look-ahead.service
[Unit]
Description=ROS2 Look-Ahead Path Following Controller
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=kikai
ExecStart=/home/kikai/ros2_ws/deploy/scripts/launch_look_ahead.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
